# deploy/systemd/pi-detector.service
[Unit]
Description=Pi5 Defect Detector
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
User=pi
WorkingDirectory=/home/pi/vision-iot-pi5/apps/pi_detector
ExecStart=/home/pi/vision/bin/python main.py
Restart=always
RestartSec=2
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
MemoryDenyWriteExecute=true
LockPersonality=true
ReadWritePaths=/home/pi/vision-iot-pi5
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
